# deployment/systemd/lightning-telegram-monitor.service
[Unit]
Description=Lightning Node Telegram Monitor
After=network.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory={{WORKING_DIRECTORY}}
ExecStart={{WORKING_DIRECTORY}}/bin/telegram-monitor
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths={{WORKING_DIRECTORY}}/data
ProtectHome=read-only

[Install]
WantedBy=default.target