# deployment/systemd/webhook-deployer.service.example
# Copy-Paste Ready Webhook Deployer Service
# Usage: 
# 1. Copy this file: sudo cp deployment/systemd/webhook-deployer.service.example /etc/systemd/system/webhook-deployer.service
# 2. Edit paths to match your setup: sudo nano /etc/systemd/system/webhook-deployer.service
# 3. Reload and start: sudo systemctl daemon-reload && sudo systemctl start webhook-deployer

[Unit]
Description=Lightning Node Tools Webhook Deployer
Documentation=https://github.com/brewgator/lightning-node-tools
After=network.target
Wants=network.target

[Service]
Type=simple
User=YOUR_USERNAME
Group=YOUR_USERNAME
WorkingDirectory=$HOME/lightning-node-tools
ExecStart=$HOME/lightning-node-tools/bin/webhook-deployer --port 9000

# Environment variables
EnvironmentFile=-$HOME/lightning-node-tools/.env

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=$HOME/lightning-node-tools/data $HOME/lightning-node-tools/logs

# Restart policy
Restart=always
RestartSec=5
StartLimitInterval=60
StartLimitBurst=3

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=webhook-deployer

[Install]
WantedBy=multi-user.target