# deployment/systemd/webhook-deployer.service
[Unit]
Description=Lightning Node Tools Webhook Deployer
Documentation=https://github.com/brewgator/lightning-node-tools
After=network.target
Wants=network.target

[Service]
Type=simple
User=__USER__
Group=__GROUP__
WorkingDirectory=__REPO_PATH__
ExecStart=__REPO_PATH__/bin/webhook-deployer \
  --port=9000 \
  --repo=__REPO_PATH__ \
  --branch=main \
  --script=./scripts/auto-deploy.sh

# Environment variables
Environment=WEBHOOK_SECRET_FILE=__REPO_PATH__/secrets/webhook.secret
EnvironmentFile=-__REPO_PATH__/.env

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=__REPO_PATH__ /var/log /tmp

# Restart policy
Restart=always
RestartSec=5
StartLimitInterval=60
StartLimitBurst=3

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=webhook-deployer

[Install]
WantedBy=multi-user.target