# deployment/systemd/webhook-deployer.service
[Unit]
Description=Lightning Node Tools Webhook Deployer
Documentation=https://github.com/brewgator/lightning-node-tools
After=network.target
Wants=network.target

[Service]
Type=simple
WorkingDirectory={{WORKING_DIRECTORY}}
ExecStart={{WORKING_DIRECTORY}}/bin/webhook-deployer --port 9000

# Environment variables
EnvironmentFile=-{{WORKING_DIRECTORY}}/.env

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{WORKING_DIRECTORY}}/data {{WORKING_DIRECTORY}}/logs

# Restart policy
Restart=always
RestartSec=5
StartLimitInterval=60
StartLimitBurst=3

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=webhook-deployer

[Install]
WantedBy=multi-user.target