# systemd/bitcoin-dashboard-collector.service
[Unit]
Description=Bitcoin Portfolio Data Collector
After=network.target bitcoin-dashboard-api.service
Wants=network-online.target

[Service]
Type=simple
User={{USER}}
Group={{USER}}
WorkingDirectory={{WORKING_DIRECTORY}}
ExecStart={{WORKING_DIRECTORY}}/bin/dashboard-collector
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths={{WORKING_DIRECTORY}}/data
ProtectHome=read-only

[Install]
WantedBy=multi-user.target